# deployments/linode/SUMSarizer.conf
[program:SUMSarizer]
command = bash /home/www/SUMSarizer/deployments/linode/start.sh
directory = /home/www/SUMSarizer
user = www

[program:SUMSarizer_worker]
command = bash /home/www/SUMSarizer/deployments/linode/start_worker.sh
directory = /home/www/SUMSarizer
user = www